# packaging/debian/postrm
#!/bin/bash
# UCM Post-Removal Script
# Handles cleanup after package removal

set -e

case "$1" in
    remove)
        # Stop service
        systemctl stop ucm.service >/dev/null 2>&1 || true
        systemctl disable ucm.service >/dev/null 2>&1 || true
        
        echo "UCM service stopped and disabled"
        echo "Data preserved in /var/lib/ucm"
        echo "To remove data: rm -rf /var/lib/ucm"
        ;;
        
    purge)
        # Remove data directories
        rm -rf /var/lib/ucm
        rm -rf /etc/ucm
        rm -rf /var/log/ucm
        
        # Remove user and group
        if getent passwd ucm >/dev/null 2>&1; then
            deluser --system ucm >/dev/null 2>&1 || true
        fi
        
        if getent group ucm >/dev/null 2>&1; then
            delgroup --system ucm >/dev/null 2>&1 || true
        fi
        
        echo "UCM completely removed"
        ;;
        
    upgrade|failed-upgrade|abort-install|abort-upgrade|disappear)
        ;;
        
    *)
        echo "postrm called with unknown argument \`$1'" >&2
        exit 1
        ;;
esac

exit 0
