# packaging/systemd/ucm-watcher.service
[Unit]
Description=UCM Service Action Handler (restart, update)
Documentation=https://github.com/NeySlim/ultimate-ca-manager
After=network.target

[Service]
Type=oneshot
ExecStart=/opt/ucm/scripts/ucm-watcher.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ucm-watcher
