# packaging/systemd/ucm-watcher.path
[Unit]
Description=Watch for UCM service actions (restart, update)
Documentation=https://github.com/NeySlim/ultimate-ca-manager

[Path]
PathExists=/opt/ucm/data/.restart_requested
PathExists=/opt/ucm/data/.update_pending
MakeDirectory=yes

[Install]
WantedBy=multi-user.target
