# packaging/rpm/ucm.service
[Unit]
Description=Ultimate CA Manager (UCM)
Documentation=https://github.com/NeySlim/ultimate-ca-manager
After=network.target

[Service]
Type=exec
User=ucm
Group=ucm
WorkingDirectory=/opt/ucm/backend

# Environment
EnvironmentFile=-/etc/ucm/ucm.env

# Start gunicorn directly using config file
ExecStart=/opt/ucm/venv/bin/gunicorn -c gunicorn_config.py wsgi:app

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ucm/data /var/log/ucm /etc/ucm
ProtectKernelTunables=true
ProtectControlGroups=true
RestrictRealtime=true
RestrictNamespaces=true
LockPersonality=true
MemoryDenyWriteExecute=false
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Resource limits
LimitNOFILE=65536
LimitNPROC=512

# Restart policy
Restart=always
RestartSec=10s

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30s

[Install]
WantedBy=multi-user.target
